# bot/sodexscan-bot.service
[Unit]
Description=SodexScan Discord Ticket Bot
After=network.target

[Service]
Type=simple
User=sodexscan
Group=sodexscan
WorkingDirectory=/opt/sodexscan-bot
ExecStart=/opt/sodexscan-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
EnvironmentFile=/opt/sodexscan-bot/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sodexscan-bot

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sodexscan-bot

[Install]
WantedBy=multi-user.target
